# root/etc/services.d/collector/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

umask "${UMASK}"

redirect_cmd() {
    if [[ "${DEBUG}" == "yes" ]]; then
        "$@"
    else
        "$@" > /dev/null 2>&1
    fi
}

s6-svwait ../web && sleep 1

while true; do
    # shellcheck disable=SC2086
    redirect_cmd exec s6-setuidgid hotio "${APP_DIR}/scrutiny-collector-metrics" run --api-endpoint "${API_ENDPOINT}" ${ARGS}
    sleep "${INTERVAL}"
done
